# rc.local
#!/bin/sh -e
#
# rc.local

python /home/pi/displayIP.py

exit 0
